# scraper-manzana/crontab
*/30 * * * * /usr/local/bin/run-scraper-once
