# scraper-py/crontab
*/1 * * * * poetry run python3 __init__.py cron